# systemd/warmoffice-off.timer
[Timer]
OnCalendar=23:00

[Install]
WantedBy=timers.targets
