# systemd/warmoffice-on.timer
[Timer]
OnCalendar=5:30
Unit=warmoffice.service

[Install]
WantedBy=timers.targets
