# systemd/warmoffice-off.service
[Unit]
Conflicts=warmoffice.service

[Service]
Type = oneshot
ExecStart = /bin/echo "stopping warmoffice"
