# systemd/warmoffice-prewarm.timer
[Timer]
OnCalendar=5:30

[Install]
WantedBy=timers.targets
